# systemd/openqa-jenkins.timer
[Unit]
Description=Grabs new build from jenkins server

[Timer]
OnCalendar=Hourly
Persistent=true

[Install]
WantedBy=timers.target
